# config/bayonne.cron
#!/bin/sh
if test -p /var/run/bayonne/control ; then
        bayonne period 60 ; fi

